# bot.conf
[program:max_way]
directory=/max-way
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/max-way/",USER="root"

stderr_logfile=max_way_error.log
stdout_logfile=max_way_log.log